# data/supervisor/sshd.conf
[program:sshd]
process_name=%(program_name)s
command = /usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/sshd/sshd.log

